# files/launch/avoutputd.conf.in

description "@WEBOS_PROJECT_SUMMARY@"


start on @START_ON_SIGNAL@

respawn

# Comment this line out to suppress logs on the console
#console output

exec @WEBOS_INSTALL_SBINDIR@/avoutputd
